# backend-py/deploy/uni-api-py.service
[Unit]
Description=University Portal Python API (FastAPI)
After=network.target postgresql.service redis-server.service

[Service]
Type=simple
User=root
WorkingDirectory=/root/University-and-Course-data/backend-py
EnvironmentFile=/root/University-and-Course-data/backend-py/.env
ExecStart=/root/University-and-Course-data/backend-py/venv/bin/gunicorn \
    app.main:app -w 4 -k uvicorn.workers.UvicornWorker -b 127.0.0.1:8000 --timeout 120 \
    --access-logfile - --error-logfile -
Restart=always
RestartSec=5
StandardOutput=journal
StandardError=journal

[Install]
WantedBy=multi-user.target
